# solvable/supervisor/webservice.conf
[program:webservice]
user=user
directory=/home/user/webservice
command=python3 server.py
